# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx-error.log

[program:api-server]
command=/usr/local/bin/api-server
autostart=true
autorestart=true
stdout_logfile=/var/log/api-server.log
stderr_logfile=/var/log/api-server-error.log
environment=PORT="8081"

[program:postfix]
command=postfix start-fg
autostart=true
autorestart=true
stdout_logfile=/var/log/postfix.log
stderr_logfile=/var/log/postfix-error.log

[program:sync-domains]
command=/usr/local/bin/sync-domains.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
